# debian/mycroft.service
[Unit]
Description=Mycroft Simple

[Service]
User=mycroft
ExecStart=/usr/bin/mycroft -wse
Restart=always

[Install]
WantedBy=multi-user.target
