# debian/prerm
#!/usr/bin/env sh

echo 'Stopping service...'
systemctl stop mycroft.service
systemctl disable mycroft.service

#DEBHELPER#

